# base-files/files/etc/uci-defaults/05_savecfg_apply
#!/bin/sh

source /lib/functions/savecfg.sh

apply_config_after_overlay_rebuild fw_redirect fw_parental wifi

exit 0
