# base-files/files/etc/uci-defaults/99-missing-sections
#!/bin/sh

uci -q get firewall.settings >/dev/null || {

uci -q batch <<-EOT
	add firewall settings
	rename firewall.@settings[-1]=settings
	set firewall.settings.disabled=0
	commit firewall
EOT

}

