# base-files/files/etc/init.d/console_log
#!/bin/sh /etc/rc.common

START=99
boot() {
    # Disable loggin to serial port.
    echo 0 >/proc/sys/kernel/printk
}
